# initscript/Gentoo/mbm-gpsd.in
#!/sbin/runscript
# Copyright (c) 2008 Saleem Abdulrasool <compnerd@compnerd.org>
# Distributed under the terms of the GNU General Purpose License v2
# $Header: $

PIDFILE=/var/run/mbm-gpsd.pid
DAEMON=/usr/sbin/mbm-gpsd

depend() {
   need udev
}

start() {
   ebegin "Starting mbm-gpsd"
   start-stop-daemon --start --quiet --pidfile ${PIDFILE} \
      --exec ${DAEMON}
   eend $?
}

stop() {
   ebegin "Stopping mbm-gpsd"
   start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
   eend $?
}

